# mpcserver.service
[Unit]
Description=mpcserver 
After=network.target

[Service]
ExecStart=/usr/local/bin/mpcserver

[Install]
WantedBy=multi-user.target
